# deploy/systemd/hexo-server.service
[Unit]
Description=Hexo Server for JBlog
After=network.target

[Service]
Type=simple
User=jeb
Group=jeb
WorkingDirectory=/home/jeb/repos/JBlog
Environment=NODE_ENV=production
ExecStart=/home/jeb/.nvm/versions/node/v24.14.0/bin/node /home/jeb/repos/JBlog/node_modules/hexo/bin/hexo server --host 0.0.0.0 --port 4000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
